# resources/nginx/beacon.service
[Unit]
Description=Chia Gaming beacon
After=network.target

[Service]
ExecStart=/app/beacon.sh @TRACKER@
